# setup/etcd-snapshot.service
[Unit]
Description=get snapshot of etcd

[Service]
Type=oneshot
ExecStart=/usr/local/bin/etcd-snapshot.sh