# setup/etcd-snapshot.timer
[Unit]
Description=get snapshot of etcd

[Timer]
OnCalendar=hourly
RandomizedDelaySec=1min
Persistent=false

[Install]
WantedBy=timers.target
