# setup/etcd-backup.timer
[Unit]
Description=get snapshot of etcd
PartOf=etcd-container.service

[Timer]
OnCalendar=hourly
RandomizedDelaySec=1min
Persistent=true
